# deploy/systemd/viralforge-secure-bot.service
[Unit]
Description=ViralForge secure Telegram command bot
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/viralforge
ExecStart=/usr/bin/docker compose -f /opt/viralforge/deploy/docker-compose.yml run --rm viralforge viralforge secure-bot
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
